# etc/rpt.service
[Unit]
Description=RPT Thermostat Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 -m src
WorkingDirectory=/home/pi/rpt
StandardOutput=inherit
StandardError=inherit
Restart=always
User=pi

[Install]
WantedBy=multi-user.target

# future work
# https://stackoverflow.com/questions/38986057/how-to-set-image-name-in-dockerfile
